# etc/rc.local
#!/bin/bash
# Backup and restore apps. Works nicely with timeshift backup.
mkdir ~/.apt-mark

apt-mark showauto > ~/.apt-mark/pkgs_auto.lst
apt-mark showmanual > ~/.apt-mark/pkgs_manual.lst
